# launchd/ai.kis.swing.eod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.kis.swing.eod</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/johyeonchang/services/kis-swing-bot/.venv/bin/python3</string>
    <string>/Users/johyeonchang/services/kis-swing-bot/src/scripts/run_eod_report.py</string>
  </array>

  <!-- 매일 15:35 실행 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>15</integer>
    <key>Minute</key>
    <integer>35</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/johyeonchang/services/kis-swing-bot</string>

  <key>StandardOutPath</key>
  <string>/Users/johyeonchang/services/kis-swing-bot/logs/eod_report.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/johyeonchang/services/kis-swing-bot/logs/eod_report.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>CLAUDE_BIN</key>
    <string>/Users/johyeonchang/.local/bin/claude</string>
    <key>PYTHONPATH</key>
    <string>/Users/johyeonchang/services/kis-swing-bot</string>
  </dict>
</dict>
</plist>
